# deploy/systemd/btcticker.service
[Unit]
Description=Bitcoin Price Ticker
After=network-online.target
Wants=network-online.target

[Service]
# Adjust User, WorkingDirectory, and ExecStart to match your install location.
User=pi
WorkingDirectory=/home/pi/raspberry/btc-ticker
ExecStart=/home/pi/raspberry/btc-ticker/.venv/bin/python -m epaper
Restart=on-failure
RestartSec=10
Type=notify
WatchdogSec=60
NotifyAccess=main

[Install]
WantedBy=multi-user.target
